# unit_files/weather-database-update.service
[Unit]
Description=Update the database with current weather measurements

[Service]
Type=oneshot
ExecStart=~/weather-pi/update_database
User=pi
