# unit_files/weather-database-update.service
[Unit]
Description=Update the database with current weather measurements

[Service]
Type=simple
ExecStart=/home/pi/weather-pi/main.py
User=pi
